# openwrt-18.06.8/v2ray.openwrt
#!/bin/sh /etc/rc.common

START=99

start() {
  start-stop-daemon -q -S -b -x /etc/v2ray/v2ray
}

stop() {
  start-stop-daemon -q -K -x /etc/v2ray/v2ray
}
